# examples/systemd/autoscaler_inputs_alertmanager.service
[Unit]
Description=SakuraCloud AutoScaler AlertManager Inputs

[Service]
User=autoscaler
EnvironmentFile=/etc/autoscaler/inputs_alertmanager.config
ExecStart=/usr/local/sbin/autoscaler inputs alertmanager $OPTIONS
Restart=always

[Install]
WantedBy=multi-user.target